# src/resources/supervisor/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd
numprocs=1
autostart=true
autorestart=true